# dockerfiles/crontab/create_file.sh
#!/bin/sh
echo "Hello world!" >> /home/test.txt